# mycat2/base/supervisord.conf
[supervisord]
nodaemon=true


[program:mysqld]
command=/entrypoint.sh mysqld



[program:mycat]
command=/usr/local/mycat2/bin/mycat start